# scripts/systemd/blackstory-overnight-enrichment.timer
# Nightly hybrid enrichment timer (Corsair user systemd).
# Default: 02:15 local — after daily SearXNG discovery (04:30 was web-search only).
[Unit]
Description=Timer for BlackStory overnight hybrid enrichment

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
